# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:rclone-daemon]
command=/usr/bin/rclone rcd --rc-addr 127.0.0.1:5572 --rc-no-auth --log-level INFO
autostart=true
autorestart=true
priority=1
stderr_logfile=/app/logs/rclone.log
stdout_logfile=/app/logs/rclone.log

[program:backend]
command=/app/server
autostart=true
autorestart=true
priority=2
stderr_logfile=/app/logs/backend.log
stdout_logfile=/app/logs/backend.log
environment=RCLONE_MANAGER_DATA_DIR="/app/data",RCLONE_MANAGER_LOG_DIR="/app/logs",RCLONE_MANAGER_PORT="7070",RCLONE_CONFIG="/root/.config/rclone/rclone.conf",TZ="Asia/Shanghai"

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
priority=3
stderr_logfile=/app/logs/nginx.log
stdout_logfile=/app/logs/nginx.log
